# files/etc/uci-defaults/wifi_enable_config 
# uci set wireless.radio1.disabled=0
# uci commit